# querierd/files/querierd.init
#!/bin/sh /etc/rc.common

START=98
USE_PROCD=1

PYTHON=/usr/bin/python
CONFIG=/etc/querierd

test -e $CONFIG || exit 0

start() {
	$PYTHON -m querier.daemon start
}

stop() {
	$PYTHON -m querier.daemon stop
}

restart() {
	$PYTHON -m querier.daemon restart
}
